# docker/php/supervisord.conf

[unix_http_server]
file=/var/run/supervisor.sock   ; (шлях до сокет файлу)
chmod=0700                      ; права доступу

[supervisord]
logfile=/var/log/supervisord.log ; головний лог-файл
pidfile=/var/run/supervisord.pid ; PID-файл
childlogdir=/var/log/            ; директорія для логів процесів
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; використання того ж сокету

[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
priority=1

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/local/bin/php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=90
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log


[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
priority=2